# s6files/services.d/wise/run
#!/bin/sh
dir=/wiseService
WISE_CONF="${WISE_CONF:-$dir/etc/wise.ini}"
cd $dir
NODE_ENV="production" exec node wiseService.js $WISE_OPTION -c $WISE_CONF
